# config/supervisord-php-fpm.conf
[program:php-fpm]
command=php-fpm --nodaemonize
priority=100
autorestart=true
