# config/supervisord-nginx.conf
[program:nginx]
command=nginx -g 'daemon off;'
priority=200
autorestart=true
